# debian/wazigate.wazigate.service
[Unit]
Description=Wazigate Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/lib/wazigate
ExecStart=/usr/bin/wazigate


[Install]
WantedBy=multi-user.target